# youtube_history_sync/com.youtube.history.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.youtube.history.sync</string>
    <key>ProgramArguments</key>
    <array>
        <!-- Replace with the absolute path to your Python interpreter (e.g., from your uv or miniforge env) -->
        <string>/Users/dhruvanand/miniforge3/bin/python3</string>
        <!-- Replace with the absolute path to the script you want to run -->
        <string>/Users/dhruvanand/Code/agent-rules-sync-standalone/youtube_history_sync/chrome_history.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/dhruvanand/youtube_history_sync_out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dhruvanand/youtube_history_sync_err.log</string>
</dict>
</plist>
